# supervisor.ini
[program:aiPy]
command=/data/www/aiPy/.venv/bin/python3.9 -m gunicorn  -c /data/www/aiPy/gu.py main.wsgi:application
directory=/data/www/aiPy/
autorestart=true
startsecs=3
startretries=3
stdout_logfile=/data/www/aiPy/logs/aiPy.supervisor.out.log
stderr_logfile=/data/www/aiPy/logs/aiPy.supervisor.err.log
stdout_logfile_maxbytes=2MB
stderr_logfile_maxbytes=2MB
user=root
numprocs=1
program_name=%(program_name)s_%(process_num)02d